# rootfs/etc/s6-overlay/s6-rc.d/airstats/run
#!/command/with-contenv sh
cd /app/core
exec s6wrap --prepend=airstats --timestamps --args /app/core/airstats
